# deploy/piassistant-kiosk.service
[Unit]
Description=PiAssistant Kiosk Display
After=piassistant.service
Wants=piassistant.service

[Service]
Type=simple
ExecStart=/usr/bin/cage -s -- /usr/bin/chromium --kiosk --noerrdialogs --disable-infobars --no-first-run --enable-features=OverlayScrollbar --disable-translate http://localhost:8000
User=%i
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=WLR_BACKENDS=drm
Environment=SEATD_VTBOUND=0
Restart=always
RestartSec=5

# TTY/seat access
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
PAMName=login
UtmpIdentifier=tty1
UtmpMode=user

# Ensure proper permissions
SupplementaryGroups=video render input

[Install]
WantedBy=multi-user.target
